# serverpulse-agent.service
[Unit]
Description=ServerPulse Monitoring Agent
Documentation=https://docs.serverpulse.com
After=network.target
Wants=network.target

[Service]
Type=simple
User=serverpulse
Group=serverpulse
WorkingDirectory=/opt/serverpulse-agent
Environment=PATH=/opt/serverpulse-agent/venv/bin
ExecStart=/opt/serverpulse-agent/venv/bin/python /opt/serverpulse-agent/serverpulse_agent.py /etc/serverpulse-agent/config.yml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
TimeoutStopSec=20

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log /etc/serverpulse-agent
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=8192
MemoryHigh=256M
MemoryMax=512M

[Install]
WantedBy=multi-user.target
